# playback/greeting.service
[Unit]
Description=Daily Greeting Flask Receiver
After=network.target

[Service]
Type=simple
User=oscar
WorkingDirectory=BASEDIR_PLACEHOLDER/playback
ExecStart=BASEDIR_PLACEHOLDER/venv/bin/python BASEDIR_PLACEHOLDER/playback/receive_greeting.py
Restart=always
RestartSec=10

# Logging
StandardOutput=append:BASEDIR_PLACEHOLDER/data/receiver.log
StandardError=append:BASEDIR_PLACEHOLDER/data/receiver.log

[Install]
WantedBy=multi-user.target
